# deploy/google-vm/take5-planner.service
[Unit]
Description=Take-5 Cabin Lake Planner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=take5
Group=take5
WorkingDirectory=/opt/take5-planner
EnvironmentFile=/etc/take5-planner.env
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /opt/take5-planner/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
